# filesystem/etc/systemd/system/archlinux-update.timer
[Unit]
Description=Arch Linux update

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
